# acpi-sleep-allow-usb-wakeup.service
[Unit]
Description=Enable ACPI Wakeup

[Service]
ExecStart=/bin/bash /usr/local/bin/acpi-sleep --enable
RemainAfterExit=true
Type=oneshot

[Install]
WantedBy=multi-user.target